# conf/osnmpd.in
#!/sbin/openrc-run
# Copyright 1999-2017 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Id$

name="SNMP agent"
pidfile="@SNMP_RUN@/osnmpd.pid"
command="@SNMP_BIN@/osnmpd"
command_args="${EXTRA_OPTS}"
start_stop_daemon_args="--pidfile ${pidfile}"

depend() {
	need localmount net
	use logger
}

